# files/usr/lib/systemd/system/hibernate-swap-resume.service
[Unit]
Description=Disable swap after resuming from hibernation
After=hibernate.target

[Service]
User=root
Type=oneshot
ExecStart=/usr/sbin/swapoff /var/swap/swapfile

[Install]
WantedBy=hibernate.target